# install/etc/systemd/system/mythtv-backend.service.d/override.conf
[Service]
UMask=0002
ExecStartPre=-/opt/mythtv/bin/mythstartup.sh
StartLimitBurst=5

